# ops/systemd/subsearch-celery.service
[Unit]
Description=Sub Search Celery Worker
After=network.target postgresql.service redis.service
Wants=network-online.target

[Service]
Type=simple
User=subsearch
Group=subsearch
WorkingDirectory=/opt/subsearch
EnvironmentFile=-/etc/subsearch.env
ExecStart=/opt/subsearch/.venv/bin/celery -A reddit_analyzer worker --loglevel=info --concurrency=1
Restart=on-failure
RestartSec=5
PrivateTmp=true
ProtectSystem=full
ProtectHome=read-only
NoNewPrivileges=true
LimitNOFILE=4096

[Install]
WantedBy=multi-user.target
