# ops/systemd/subsearch.service
[Unit]
Description=Sub Search (Django) via Gunicorn
After=network.target postgresql.service redis.service
Wants=network-online.target

[Service]
Type=simple
User=subsearch
Group=subsearch
WorkingDirectory=/opt/subsearch
EnvironmentFile=-/etc/subsearch.env
ExecStart=/opt/subsearch/venv/bin/gunicorn reddit_analyzer.wsgi:application --workers 4 --bind 127.0.0.1:8000
Restart=on-failure
RestartSec=5
PrivateTmp=true
ProtectSystem=full
ProtectHome=read-only
NoNewPrivileges=true
LimitNOFILE=4096

[Install]
WantedBy=multi-user.target
